# backend/app/templates/systemd/schema-registry.service.j2
[Unit]
Description=Apicurio Schema Registry (Tantor-managed)
After=network.target kafka.service
Wants=kafka.service

[Service]
Type=simple
User=apicurio
Group=apicurio
Environment="JAVA_HOME={{ java_home }}"
# v1.4.0 — Apicurio 3.x stopped honoring QUARKUS_CONFIG_LOCATIONS
# at startup, so we pass the critical config keys as -D system
# properties on the command line. The application.properties file
# stays around for human inspection / future use.
Environment="QUARKUS_CONFIG_LOCATIONS={{ config_path }}"
{% set sr_java_opts = [] %}
{% if heap_opts %}{% set _ = sr_java_opts.append(heap_opts) %}{% endif %}
{% if jvm_performance_opts %}{% set _ = sr_java_opts.append(jvm_performance_opts) %}{% endif %}
{% if gc_logging_enabled %}{% set _ = sr_java_opts.append("-Xlog:gc*:file=/var/log/kafka/gc-schema-registry.log:time,level,tags:filecount=10,filesize=10M") %}{% endif %}
{% if jmx_port %}{% set _ = sr_java_opts.append("-Dcom.sun.management.jmxremote -Dcom.sun.management.jmxremote.port=" ~ jmx_port ~ " -Dcom.sun.management.jmxremote.rmi.port=" ~ jmx_port ~ " -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false") %}{% endif %}
{% if sr_java_opts %}
Environment="JAVA_OPTS={{ sr_java_opts | join(' ') }}"
{% endif %}
# Apicurio 3.1+ ships quarkus-run.jar at the tarball root after our
# strip-components extract. The legacy quarkus-app/ subdir layout was
# removed in 3.0. Pick whichever exists.
ExecStart=/bin/sh -c 'JAR=/opt/apicurio/quarkus-app/quarkus-run.jar; [ -f "$JAR" ] || JAR=/opt/apicurio/quarkus-run.jar; exec {{ java_home }}/bin/java $JAVA_OPTS -Dquarkus.http.port={{ schema_registry_port | default(8085) }} -Dapicurio.storage.kind=kafkasql -Dapicurio.kafkasql.bootstrap.servers={{ bootstrap_servers | default("127.0.0.1:9092") }} -jar $JAR'
Restart=always
RestartSec=5
# v1.4.1 — Quarkus cold-start + kafkasql topic creation can exceed
# systemd's default 90s start window on small VMs. 300s avoids spurious
# "service failed" while still failing fast on a real bug.
TimeoutStartSec=300
LimitNOFILE=65536
{% if cpu_quota %}
CPUQuota={{ cpu_quota }}
{% endif %}
{% if memory_max %}
MemoryMax={{ memory_max }}
{% endif %}

[Install]
WantedBy=multi-user.target
